# service/backup.service
[Unit]
Description=Service for borgbackup etc files

[Service]
Type=oneshot
ExecStart=/vagrant/scripts/borg.sh

[Install]
WantedBy=multi-user.target
